# Capas_personalizadas/meta-server/recipes-server/server/files/server.service
[Unit]
Description=Servidor personalizado
After=network.target

[Service]
ExecStart=/usr/bin/server
Restart=always

[Install]
WantedBy=multi-user.target
